# user-init-scripts/profile.d/acpp-path.sh
export PATH=/usr/local/acpp/bin:"$PATH"
